# roles/frontend/templates/npmstart.service
[Unit]
Description=npm service

[Service]

WorkingDirectory=/root
ExecStart=npm start
Restart=on-failure

[Install]
WantedBy=multi-user.target